# L4-SD-WAN/Otros/2disable-offloads.service
[Unit]
Description=Disable NIC offloads inside LXC
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/disable-offloads.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
